# etc/systemd/system/agostle.socket
[Unit]
Description=Agostle listening socket

[Socket]
Service=agostle@1.service
ListenStream=9500
NoDelay=true

[Install]
WantedBy=local.target
